# config/vibeswaf.service
[Unit]
Description=Vibeswaf WAF Backend Service
After=network.target postgresql.service clickhouse-server.service
Wants=postgresql.service clickhouse-server.service

[Service]
Type=simple
User=vibeswaf
Group=vibeswaf
WorkingDirectory=/opt/vibeswaf
ExecStart=/opt/vibeswaf/vibeswaf
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/vibeswaf/.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# acme.sh writes account state, renewal config, and per-domain cert storage
# under /opt/vibeswaf/.acme.sh, so it must be writable under ProtectSystem=strict.
ReadWritePaths=/opt/vibeswaf/logs /opt/vibeswaf/data /opt/vibeswaf/.acme.sh /opt/certs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
